# home/.config/autostart/akregator.desktop
[Desktop Entry]
Version=1.0
Name=Akregator
Exec=akregator --hide-mainwindow -- %u
Terminal=false
Type=Application
