# home/.config/autostart/autostart.desktop
[Desktop Entry]
Exec=autostart
Icon=dialog-scripts
Name=autostart
Path=
Type=Application
X-KDE-AutostartScript=true
