# home/.config/autostart/autostart.desktop
[Desktop Entry]
Exec=/home/matklad/.local/bin/autostart
Icon=dialog-scripts
Name=autostart
Path=
Type=Application
X-KDE-AutostartScript=true
